# systemd/jellyweb.service
[Unit]
Description=JellyWatch Web UI Server
Documentation=https://github.com/Nomadcxx/jellywatch
After=network.target jellywatchd.service
Wants=jellywatchd.service

[Service]
Type=simple
User=root
Group=root
Environment=SUDO_USER=nomadx
ExecStart=/usr/local/bin/jellyweb --host 0.0.0.0 --port 5522
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jellyweb

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome left unset: jellyweb needs RW on the user's config dir
# (config.toml + lock files + media.db live under ~/.config/jellywatch).
# ReadWritePaths over a ProtectHome=read-only mount has proven unreliable
# in practice (config saves return EROFS), so we leave /home accessible.
ReadWritePaths=/home
PrivateTmp=true

[Install]
WantedBy=multi-user.target
