# systemd/plex2jellyfin-web.service
[Unit]
Description=Plex2Jellyfin Web UI Server
Documentation=https://github.com/Nomadcxx/plex2jellyfin
After=network.target plex2jellyfin-daemon.service
Wants=plex2jellyfin-daemon.service

[Service]
Type=simple
User=root
Group=root
# SUDO_USER tells the web UI whose ~/.config/plex2jellyfin to read.
# Package installs must set it once via a drop-in:
#   sudo systemctl edit plex2jellyfin-web
#   [Service]
#   Environment=SUDO_USER=<your username>
ExecStart=/usr/bin/plex2jellyfin-web --host 0.0.0.0 --port 5522
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plex2jellyfin-web

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
# ProtectHome left unset: plex2jellyfin-web needs RW on the user's config dir
# (config.toml + lock files + media.db live under ~/.config/plex2jellyfin).
# Plex2Jellyfin web also executes user-triggered library maintenance actions such as
# consolidation, so non-system media mounts must remain writable according to
# their normal filesystem permissions rather than a hardcoded allow-list.
PrivateTmp=true

[Install]
WantedBy=multi-user.target
