# systemd/plex2jellyfin-daemon.service
[Unit]
Description=Plex2Jellyfin Media File Watcher
After=network.target

[Service]
Type=simple
User=root
Group=root
# SUDO_USER tells the daemon whose ~/.config/plex2jellyfin to read.
# Package installs must set it once via a drop-in:
#   sudo systemctl edit plex2jellyfin-daemon
#   [Service]
#   Environment=SUDO_USER=<your username>
ExecStart=/usr/bin/plex2jellyfin-daemon
Restart=on-failure
RestartSec=10

# Security settings
# NoNewPrivileges disabled - daemon needs to chown files
# NoNewPrivileges=true
PrivateTmp=true

# Restrict capabilities to minimum needed
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target
