# systemd/user/jellywatch-postmortem.timer
[Unit]
Description=Run JellyWatch postmortem every four days

[Timer]
OnBootSec=20m
OnUnitActiveSec=96h
Persistent=true
Unit=jellywatch-postmortem.service

[Install]
WantedBy=timers.target
