# systemd/user/plex2jellyfin-postmortem.service
[Unit]
Description=Plex2Jellyfin scheduled postmortem collector

[Service]
Type=oneshot
WorkingDirectory=/home/nomadx/Documents/plex2jellyfin
ExecStart=/usr/local/bin/plex2jellyfin postmortem collect --since 96h
ExecStartPost=/home/nomadx/Documents/plex2jellyfin/scripts/plex2jellyfin-postmortem-terminal.sh
